# backup_scripts/mysql-backup.timer
[Unit]
Descripton=Run MySQL Backup Service Weekly

[Timer]
OnCalendar=weekly
Persistent=true

[Install]
WantedBy=timers.target
